# configs/groups/dev/environment/supervisor/supervisor.ini
; Запустить supervisor в основном (не фоновом) режиме
[supervisord]
nodaemon=true

; Запустить бинарник nginx
; - user: запускать от имени пользователя root
; - autostart: запускать автоматически при старте supervisor
[program:nginx]

; Команда, запускать ли при старте, от имени какого пользователя
command = /usr/sbin/nginx
user = root
autostart = true
autorestart = false

; Где вести лог ошибок сервиса
redirect_stderr=false
stderr_logfile=<path_to_the_project>/other/logs/nginx
stderr_logfile_maxbytes=1MB

; Запустить бинарник php-fpm7.0
[program:php-fpm]

; Команда, запускать ли при старте, от имени какого пользователя
command = /usr/sbin/php-fpm7.0 -F
user = root
autostart = true
autorestart = false

; Где вести лог ошибок сервиса
redirect_stderr=false
stderr_logfile=<path_to_the_project>/other/logs/php-fpm
stderr_logfile_maxbytes=1MB

; Запустить cron
[program:cron]

; Команда, запускать ли при старте, от имени какого пользователя, нужен ли авторестарт
command = /usr/sbin/cron -f
user = root
autostart = true
autorestart = false

; Где вести лог ошибок сервиса
redirect_stderr=false
stderr_logfile=<path_to_the_project>/other/logs/cron
stderr_logfile_maxbytes=1MB

; Запустить ларавеловского демона queue:listen 
[program:queue]

; Команда, запускать ли при старте, от имени какого пользователя, нужен ли авторестарт
command = php <path_to_the_project>/project/artisan queue:listen
user = root
autostart = true
autorestart = false

; Где вести лог ошибок сервиса
redirect_stderr=false
stderr_logfile=<path_to_the_project>/other/logs/queue
stderr_logfile_maxbytes=1MB
